# supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/usr/bin/mongod --smallfiles

[program:workbench_server]
command=/usr/local/bin/workbench_server
